# files/home/config/systemd/user/mpDris2.service.d/override.conf.j2
{{ ansible_managed | comment('plain', prefix='#################################', postfix='#################################') }}
[Unit]
After=mpd.service

[Install]
RequiredBy=mpd.service
