# files/home/config/systemd/user/greenclip.service.d/override.conf.j2
{{ ansible_managed | comment('plain', prefix='#################################', postfix='#################################') }}
[Service]
RestartSec=10
Restart=always
